# root/etc/services.d/chowner/run
#!/usr/bin/with-contenv bash

exec \
    s6-setuidgid abc \
    find /data \
    -not -user abc -o \
    -not -group abc \
    -execdir chown abc:abc {} \+
