# root/etc/services.d/chowner/run
#!/usr/bin/with-contenv bash

FIND_PATHS=${FIND_PATHS:-/data}

/usr/bin/find ${FIND_PATHS} \
\( -not -user abc -o \
-not -group abc \) \
-exec ls {} \; -exec chown abc:abc {} \+

sleep 5
